# debian/fhs/etc/cron.daily/fexp-cleanup
#!/usr/bin/php5
<?php // INDENTING (emacs/vi): -*- mode:html; tab-width:2; c-basic-offset:2; intent-tabs-mode:nil; -*- ex: set tabstop=2 expandtab:
require_once 'FEXP.php';
$oFEXP = new FEXP( '/etc/php5/fexp/config.php' );
$oFEXP->doAutoCleanup();

